# debian/metrics-monitor-system/DEBIAN/prerm
#!/bin/sh
set -e

# Pre-removal script for metrics-monitor-system

# Stop and disable service if running
if [ -f /lib/systemd/system/metrics-monitor-system.service ]; then
    systemctl stop metrics-monitor-system >/dev/null 2>&1 || true
    systemctl disable metrics-monitor-system >/dev/null 2>&1 || true
    systemctl daemon-reload >/dev/null 2>&1 || true
fi

# DEBHELPER#
exit 0

